# deploy/polymarket-ingestion.service.d/limits.conf
# systemd drop-in — resource limits for polymarket-ingestion
#
# Overrides the base unit's MemoryMax (1G) — actual effective value is 512M.
# Ingestion is lightweight (HTTP fetch + DB upsert); 512M is more than
# sufficient and prevents it from competing with trading bots for RAM.
# OOMScoreAdjust=+100 makes ingestion MORE likely to be OOM-killed than
# trading bots — if memory is tight, losing an ingestion cycle is preferable
# to losing a trading process.
#
# Installed at: /etc/systemd/system/polymarket-ingestion.service.d/limits.conf
# NOT auto-installed by deploy.sh (which copies only *.service files).
# Committed to repo S235 (2026-05-31) for documentation + host-rebuild recovery.

[Service]
MemoryMax=512M
OOMScoreAdjust=100
