# deploy/polymarket-mirror.service.d/limits.conf
# systemd drop-in — resource limits for polymarket-mirror (MirrorBot)
#
# Overrides the base unit's MemoryMax (3G) — actual effective value is 2560M.
# OOMScoreAdjust=-100 makes the kernel prefer to OOM-kill other processes
# before MirrorBot (negative = less likely to be killed; -100 is moderate
# protection; lower than weather at -200 because mirror has more stable
# memory usage than the weather model stack).
#
# Installed at: /etc/systemd/system/polymarket-mirror.service.d/limits.conf
# NOT auto-installed by deploy.sh (which copies only *.service files).
# Committed to repo S235 (2026-05-31) for documentation + host-rebuild recovery.

[Service]
MemoryMax=2560M
OOMScoreAdjust=-100
