# deploy/polymarket-ai.service
[Unit]
Description=Polymarket AI V2 Trading System
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-ai-v2
Environment=PYTHONUNBUFFERED=1
Environment="PYTHONWARNINGS=ignore::UserWarning:sklearn,ignore::FutureWarning:sklearn"
EnvironmentFile=/opt/pa2-shared/.env
ExecStart=/opt/polymarket-ai-v2/venv/bin/python main.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# App has its own internal watchdog (_watchdog() in main.py) for bot liveness.
# systemd watchdog disabled — app does not send sd_notify heartbeats.
# NotifyAccess=all

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/pa2-shared/data /var/log/polymarket
PrivateTmp=true

# Resource limits
MemoryMax=6G
CPUQuota=180%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polymarket-ai

[Install]
WantedBy=multi-user.target
