# deploy/polymarket-esports-ingestion.service
[Unit]
Description=Polymarket AI V2 — EsportsBot Ingestion Service (EB splinter)
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

# Splinter contract (mirrors polymarket-esports.service.d/00-splinter.conf):
#   - Runs from EB splinter release path /opt/polymarket-ai-v2-esports
#   - Reads EB-specific ingestion env from /opt/pa2-shared/.env.esports.ingestion
#   - Filtered to EB bots ONLY (EsportsBot + EsportsBotV2) via env var
#     INGESTION_BOT_FILTER_MODE=include / INGESTION_BOT_FILTER_LIST=...
#   - Does NOT touch MirrorBot or WeatherBot rows
#   - Master polymarket-ingestion.service runs in parallel with the inverse
#     filter (mode=exclude), see EB_COORDINATION_RESOLUTION_VS_CHAIN.md
#
# Filter logic implementation lives in base_engine/data/resolution_backfill.py
# and is included from the same splinter release. Until that filter logic
# lands (in-progress), this unit MUST NOT be enabled — running it in parallel
# with master ingestion would race on RESOLUTION writes via NOT EXISTS.

[Service]
Type=simple
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-ai-v2-esports
Environment=PYTHONUNBUFFERED=1
Environment="PYTHONWARNINGS=ignore::UserWarning:sklearn,ignore::FutureWarning:sklearn"
EnvironmentFile=-/opt/pa2-shared/.env
EnvironmentFile=-/opt/pa2-shared/.env.esports.ingestion
ExecStart=/opt/pa2-esports-shared/venv/bin/python ingestion_main.py
Restart=always
RestartSec=15
TimeoutStopSec=60
StartLimitIntervalSec=300
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/pa2-shared/data /opt/pa2-shared/saved_models /var/log/polymarket
PrivateTmp=true

# EB-scoped ingestion has a smaller working set than master ingestion.
MemoryMax=768M
CPUQuota=40%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polymarket-esports-ingestion

[Install]
WantedBy=multi-user.target
