# deploy/polymarket-orderbook.timer
[Unit]
Description=Polymarket Orderbook Collector Timer

[Timer]
OnCalendar=*:*:00
Persistent=true

[Install]
WantedBy=timers.target
