# deploy/polymarket-mirror.service
[Unit]
Description=Polymarket AI V2 — MirrorBot
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
# S249: Type=notify + WatchdogSec defends against scan-loop wedges that
# the in-process watchdog cannot catch (S247-WB ⑥ pattern: bot wedged
# ~44h, NRestarts=0 because the asyncio loop froze but the process
# stayed alive). main.py sends READY=1 after startup and WATCHDOG=1
# every ~30s from the watchdog loop. 120s = 4× the watchdog interval —
# survives a couple of slow iterations, kills on a true wedge.
Type=notify
NotifyAccess=main
WatchdogSec=120s
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-ai-v2
Environment=PYTHONUNBUFFERED=1
Environment="PYTHONWARNINGS=ignore::UserWarning:sklearn,ignore::FutureWarning:sklearn"
EnvironmentFile=/opt/pa2-shared/.env
EnvironmentFile=/opt/pa2-shared/.env.mirror
ExecStart=/opt/polymarket-ai-v2/venv/bin/python main.py
Restart=always
RestartSec=10
TimeoutStopSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/pa2-shared/data /opt/pa2-shared/saved_models /var/log/polymarket
PrivateTmp=true

# Resource limits — MirrorBot heaviest (RTDS + 500 elites + open positions)
MemoryMax=3G
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polymarket-mirror

[Install]
WantedBy=multi-user.target
