# deploy/polymarket-orderbook.service
[Unit]
Description=Polymarket Orderbook Collector
After=postgresql@16-main.service

[Service]
Type=oneshot
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-ai-v2
ExecStart=/opt/pa2-shared/venv/bin/python scripts/orderbook_collector.py --once --limit 200
Environment=PYTHONPATH=/opt/polymarket-ai-v2
EnvironmentFile=/opt/pa2-shared/.env
TimeoutStartSec=90
MemoryMax=256M
