# deploy/polymarket-audit.service
[Unit]
Description=Run data integrity audit (21+ checks)
After=postgresql.service

[Service]
Type=oneshot
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-ai-v2
EnvironmentFile=/opt/polymarket-ai-v2/.env
Environment=PYTHONPATH=/opt/polymarket-ai-v2
ExecStart=/opt/pa2-shared/venv/bin/python scripts/run_audit.py
TimeoutStartSec=900

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polymarket-audit
