# deploy/supervisord.conf
; Polymarket AI V2 — Supervisor Process Management
; Alternative to systemd for environments where systemd is unavailable.
;
; Install: pip install supervisor
; Run: supervisord -c deploy/supervisord.conf
; Status: supervisorctl -c deploy/supervisord.conf status
; Restart: supervisorctl -c deploy/supervisord.conf restart all

[unix_http_server]
file=/tmp/polymarket_supervisor.sock

[supervisord]
logfile=/var/log/polymarket/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info
pidfile=/tmp/polymarket_supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/polymarket_supervisor.sock

; ── Main trading engine ──────────────────────────────────────────
[program:polymarket-ai]
command=/usr/bin/python3 main.py
directory=/opt/polymarket-ai-v2
environment=PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/log/polymarket/trading.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
priority=100

; ── Streamlit dashboard ──────────────────────────────────────────
[program:polymarket-dashboard]
command=/usr/bin/python3 -m streamlit run ui/dashboard.py --server.port 8501 --server.address 0.0.0.0
directory=/opt/polymarket-ai-v2
environment=PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/polymarket/dashboard.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
priority=200

; ── Data ingestion (periodic, managed by cron normally) ──────────
[program:polymarket-ingestion]
command=/usr/bin/python3 scripts/ingest_data.py
directory=/opt/polymarket-ai-v2
environment=PYTHONUNBUFFERED="1"
autostart=false
autorestart=false
startsecs=5
redirect_stderr=true
stdout_logfile=/var/log/polymarket/ingestion.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=2
priority=300

; ── Group all processes ──────────────────────────────────────────
[group:polymarket]
programs=polymarket-ai,polymarket-dashboard
priority=100
