# deploy/polymarket-esports.service.d/00-splinter.conf
# systemd drop-in override — EB SPLINTER (eb/main branch)
#
# This file overrides specific properties of polymarket-esports.service to point
# the bot at the splinter release path (/opt/polymarket-ai-v2-esports) instead
# of master's path (/opt/polymarket-ai-v2). It is installed to:
#   /etc/systemd/system/polymarket-esports.service.d/00-splinter.conf
#
# Standard systemd drop-in semantics: loaded ON TOP of the main service file
# at /etc/systemd/system/polymarket-esports.service, regardless of what's in
# the main file. So when master's deploy.sh re-installs its own
# polymarket-esports.service (pointing at /opt/polymarket-ai-v2), this override
# still wins — EB stays on the splinter.
#
# The empty `ExecStart=` line clears the inherited value before setting the
# new one (systemd's drop-in semantic for list-valued options).

[Service]
WorkingDirectory=/opt/polymarket-ai-v2-esports
ExecStart=
ExecStart=/opt/polymarket-ai-v2-esports/venv/bin/python main.py
