# deploy/polymarket-esports.service.d/limits.conf
# systemd drop-in — resource limits for polymarket-esports (EsportsBotV2)
#
# Overrides the base unit's MemoryMax (2G) — actual effective value is 2560M.
# EB v2 runs larger model stacks than EB v1; the higher limit prevents OOM
# during inference on match days with many concurrent predictions.
# OOMScoreAdjust=0 (neutral — no special OOM protection or exposure).
#
# Installed at: /etc/systemd/system/polymarket-esports.service.d/limits.conf
# NOT auto-installed by deploy.sh (which copies only *.service files).
# Committed to repo S235 (2026-05-31) for documentation + host-rebuild recovery.

[Service]
MemoryMax=2560M
OOMScoreAdjust=0
