# deploy/polymarket-weather.service.d/limits.conf
# systemd drop-in — resource limits for polymarket-weather (WeatherBot)
#
# MemoryMax=2G matches the base unit — no effective drift for memory.
# OOMScoreAdjust=-200 gives WeatherBot stronger OOM protection than mirror
# or esports. The weather model stack (sklearn + catboost + large numpy
# arrays) spikes memory during retrain; protecting it reduces retrain-induced
# service interruptions.
#
# Installed at: /etc/systemd/system/polymarket-weather.service.d/limits.conf
# NOT auto-installed by deploy.sh (which copies only *.service files).
# Committed to repo S235 (2026-05-31) for documentation + host-rebuild recovery.

[Service]
MemoryMax=2G
OOMScoreAdjust=-200
